# ami-images/simian-army/base/files/amazon/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.

# TO ADD HERE
